# perfana-kedro/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:jupyter]
user=%(ENV_NB_UID)s
command=start.sh jupyter lab --NotebookApp.token=''

[program:sshd]
command=/usr/sbin/sshd -D
